# crontab
*/5 * * * * sh -c '/app/hetzner_traffic_manager.sh 2>&1'