# crontab
*/5 * * * * /app/hetzner_traffic_manager.sh >> /proc/1/fd/1 2>&1
